# geonode/supervisor-app.conf
[program:app-geoserver]
directory = /home/vagrant/geonode
command = /home/vagrant/.venvs/geonode/bin/paver start_geoserver
